# rlclaw-dashboard.service
[Unit]
Description=rlclaw dashboard
After=network.target

[Service]
Type=simple
User=jacob
WorkingDirectory=/home/jacob/rlclaw
ExecStart=/usr/bin/npx tsx src/dashboard/server.ts
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
